# var/spool/cron/crontabs/man
# DO NOT EDIT THIS FILE - edit the master and reinstall.
#
0 0 * * * /usr/sbin/makewhatis
1 0 * * * test -d /usr/local/man && /usr/sbin/makewhatis /usr/local/man
