# var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
#
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin

*/3 *  * * *		/usr/lib/atrun
30  7  1 1 *		savelogs yearly
35  7  1 * *		savelogs monthly
40  7  * * 1		savelogs weekly
40  8  * * 1		locate.updatedb
45  7  * * *		savelogs
